# override/etc/services.d/00-unifi/run
#!/usr/bin/with-contenv sh
HOME=/mnt/config/home
cd /srv/unifi
HOME=/mnt/config/home exec s6-setuidgid guardian java -Xmx1024M -jar /srv/unifi/lib/ace.jar start;