# override/etc/services.d/10-unifi/run
#!/usr/bin/with-contenv sh
cd /srv/unifi
exec s6-setuidgid guardian java -Xmx1024M -jar /srv/unifi/lib/ace.jar start;